# deprecated/example.crontab
*/5 *   *   *   *    ~/office/checkstallman.sh
